# files/etc/rc.local
#!/bin/sh

/etc/init.d/gc restart
 (sleep 25 && /etc/init.d/modemmanager restart) & 
 ( procd -h /etc/hotplug-gcreboot.json> /dev/null ) & 

chmod +x /etc/init.d/auto-modem
/etc/init.d/auto-modem enable



exit 0
